# Services/Systemctl/ClamavAPI.service
[Unit]
Description=Running docker for ClamavAPI App

[Service]
ExecStart=/bin/bash ../Scripts/ClamavAPI.sh

[Install]
WantedBy=multi-user.target
