# Services/Cron/crontab
0 0 * * 0  ../Scripts/RestartServer.sh
* * * * * ../Scripts/FreeMemoryAlert.sh
